# conf/autostart.service
[Unit]
Description=autostart service
After=network.target
[Service]
ExecStart=/etc/init.d/autostart
[Install]
WantedBy=multi-user.target
