# conf/rfkillunblock.service
[Unit]
Description=rfkillunblock service
After=network.target
[Service]
ExecStart=/etc/init.d/rfkillunblock
[Install]
WantedBy=multi-user.target
